# com.smarttrade.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.smarttrade.scheduler</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/manirastegari/maniProject/SmartTrade-AI-UltimateStrategy/automated_daily_scheduler.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/manirastegari/maniProject/SmartTrade-AI-UltimateStrategy</string>
    
    <key>StandardOutPath</key>
    <string>/Users/manirastegari/maniProject/SmartTrade-AI-UltimateStrategy/scheduler_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/manirastegari/maniProject/SmartTrade-AI-UltimateStrategy/scheduler_stderr.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>SMARTTRADE_PROJECT_PATH</key>
        <string>/Users/manirastegari/maniProject/SmartTrade-AI-UltimateStrategy</string>
    </dict>
</dict>
</plist>
